# shutdown-backup-enable@.service
[Unit]
Description=Enable shutdown backup

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl enable shutdown-backup@%i
